# crontab
00 00 * * * echo "What\'s the time, Mr. Wolf?"
05 0-23/7 * * 7 perl blah.pl
00 1,3,7,12 * * * ps aux 
